# tools/launchagents/ai.openclaw.heartbeat-sender.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.heartbeat-sender</string>

  <!-- 使用 wrapper 脚本: 先检查 Bot 进程存活才发心跳 -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/blackdj/Desktop/OpenClaw Bot/packages/clawbot/scripts/heartbeat_sender.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/Users/blackdj/Desktop/OpenClaw Bot/packages/clawbot/logs/heartbeat-sender.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/blackdj/Desktop/OpenClaw Bot/packages/clawbot/logs/heartbeat-sender.log</string>
</dict>
</plist>
