# tools/launchagents/ai.openclaw.kiro-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.kiro-gateway</string>

  <!-- 内联启动命令：macOS 26+ Sandbox 策略禁止 bash 读取带 provenance 属性的脚本文件，必须把命令内联到 -c 参数中 -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd "/Users/blackdj/Desktop/OpenEverything/packages/clawbot/kiro-gateway" &amp;&amp; exec "/Users/blackdj/Desktop/OpenEverything/packages/clawbot/kiro-gateway/.venv/bin/python" "/Users/blackdj/Desktop/OpenEverything/packages/clawbot/kiro-gateway/main.py" --host 127.0.0.1 --port 18793</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>15</integer>

  <!-- 日志路径必须在 ~/Library/Logs/ 下：macOS 26+ Sandbox 策略禁止 xpcproxy 读取带 provenance 属性的项目目录下的文件 -->
  <key>StandardOutPath</key>
  <string>/Users/blackdj/Library/Logs/OpenClaw/kiro-gateway.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/blackdj/Library/Logs/OpenClaw/kiro-gateway.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/blackdj/Desktop/OpenEverything/packages/clawbot/kiro-gateway/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>LANG</key>
    <string>en_US.UTF-8</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
