# tools/launchagents/ai.openclaw.clawbot-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.clawbot-agent</string>

  <!-- 用 /bin/bash -c exec 间接调用：绕过 macOS provenance 安全检查（com.apple.provenance 扩展属性导致 launchd 直接执行带标记文件时返回 exit 78） -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>exec "/Users/blackdj/Desktop/OpenClaw Bot/packages/clawbot/.venv312/bin/python" "/Users/blackdj/Desktop/OpenClaw Bot/packages/clawbot/multi_main.py"</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/blackdj/Desktop/OpenClaw Bot/packages/clawbot</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/blackdj/Desktop/OpenClaw Bot/packages/clawbot/logs/com-clawbot-agent.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/blackdj/Desktop/OpenClaw Bot/packages/clawbot/logs/com-clawbot-agent.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/blackdj/Desktop/OpenClaw Bot/packages/clawbot/.venv312/bin:/opt/homebrew/bin:/usr/local/bin:/usr/sbin:/usr/bin:/bin:/sbin</string>
    <key>LANG</key>
    <string>en_US.UTF-8</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
