# tools/launchagents/ai.openclaw.g4f.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.g4f</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/blackdj/Desktop/OpenClaw Bot/clawbot/browser-agent/.venv/bin/python</string>
    <string>-m</string>
    <string>g4f</string>
    <string>api</string>
    <string>--port</string>
    <string>18891</string>
    <string>--g4f-api-key</string>
    <string>dummy</string>
    <string>--no-gui</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/blackdj/Desktop/OpenClaw Bot/clawbot/browser-agent</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/blackdj/Desktop/OpenClaw Bot/clawbot/logs/com-clawbot-g4f.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/blackdj/Desktop/OpenClaw Bot/clawbot/logs/com-clawbot-g4f.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/blackdj/Desktop/OpenClaw Bot/clawbot/browser-agent/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>LANG</key>
    <string>en_US.UTF-8</string>
  </dict>
</dict>
</plist>
